# api/scripts/dev.tylerkeller.aw-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd plist for aw_sync.py

  Setup:
    1. Fill in AW_SYNC_KEY and the ExecStart path below.
    2. cp dev.tylerkeller.aw-sync.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/dev.tylerkeller.aw-sync.plist

  Check status:  launchctl list dev.tylerkeller.aw-sync
  Tail logs:     tail -f /tmp/aw-sync.log
  Unload:        launchctl unload ~/Library/LaunchAgents/dev.tylerkeller.aw-sync.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.tylerkeller.aw-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <!-- absolute path to this script -->
        <string>/path/to/api/scripts/aw_sync.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AW_SYNC_KEY</key>
        <string>REPLACE_WITH_YOUR_API_KEY</string>
        <key>AW_SYNC_URL</key>
        <string>https://api.tylerkeller.dev</string>
        <!-- optional overrides:
        <key>AW_SYNC_MACHINE</key>
        <string>macbook</string>
        <key>AW_SYNC_INTERVAL</key>
        <string>300</string>
        -->
    </dict>

    <!-- start on login and keep alive if it exits -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/aw-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/aw-sync.log</string>
</dict>
</plist>
