# api/scripts/dev.tylerkeller.aw-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd agent for aw_sync.py

  SETUP
  1. Find your Python that has `requests` installed:
       which python3        (Homebrew: /opt/homebrew/bin/python3)
       python3 -c "import requests; print('ok')"

  2. Edit the two REPLACE_ME values below (python path + script path).

  3. Copy to LaunchAgents and bootstrap:
       cp dev.tylerkeller.aw-sync.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/dev.tylerkeller.aw-sync.plist

  USEFUL COMMANDS
  Status:   launchctl print gui/$(id -u)/dev.tylerkeller.aw-sync
  Logs:     tail -f /tmp/aw-sync.log
  Stop:     launchctl bootout gui/$(id -u)/dev.tylerkeller.aw-sync
  Restart:  launchctl kickstart -k gui/$(id -u)/dev.tylerkeller.aw-sync
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.tylerkeller.aw-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>REPLACE_ME_python3_path</string>
        <string>REPLACE_ME_absolute_path_to_aw_sync.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AW_SYNC_KEY</key>
        <string>REPLACE_WITH_YOUR_API_KEY</string>
        <key>AW_SYNC_URL</key>
        <string>https://api.tylerkeller.dev</string>
        <key>AW_SYNC_MACHINE</key>
        <string>personal-macbook</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/aw-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/aw-sync.log</string>
</dict>
</plist>
