# api/scripts/dev.tylerkeller.aw-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd agent for aw_sync.py

  FIRST-TIME SETUP
  1. Store the API key in Keychain (never stored in this file):
       security add-generic-password -a $USER -s aw-sync-key -w "your-key-here"

  2. Copy to LaunchAgents and load:
       cp dev.tylerkeller.aw-sync.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/dev.tylerkeller.aw-sync.plist

  3. Verify it's running:
       launchctl print gui/$(id -u)/dev.tylerkeller.aw-sync

  RELOADING AFTER CHANGES (bootstrap fails if already registered — always bootout first)
       cp dev.tylerkeller.aw-sync.plist ~/Library/LaunchAgents/
       launchctl bootout gui/$(id -u)/dev.tylerkeller.aw-sync 2>/dev/null; true
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/dev.tylerkeller.aw-sync.plist

  USEFUL COMMANDS
  Status:   launchctl print gui/$(id -u)/dev.tylerkeller.aw-sync
  Logs:     tail -f /tmp/aw-sync.log
  Stop:     launchctl bootout gui/$(id -u)/dev.tylerkeller.aw-sync
  Restart:  launchctl kickstart -k gui/$(id -u)/dev.tylerkeller.aw-sync
  Update key: security add-generic-password -U -a $USER -s aw-sync-key -w "new-key"
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.tylerkeller.aw-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tylerkeller/.local/share/aw-sync/venv/bin/python</string>
        <string>/Users/tylerkeller/Code/Personal/Apps/WebApps/tylerkeller-dev/api/scripts/aw_sync.py</string>
    </array>


    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/aw-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/aw-sync.log</string>
</dict>
</plist>
